# etc/supervisor/monerod.conf
[supervisord]
nodaemon=true
#loglevel=trace

[program:monerod]
user=monero
environment=HOME="/home/monero",USER="monero"
command=/home/monero/monero-gui/monerod --check-updates disabled --log-level 2
redirect_stderr=true
stdout_logfile=/home/monero/%(program_name)s.log
stderr_logfile=/home/monero/%(program_name)s.log
autorestart=true
